# etc/systemd/system/systemd-udevd.service.d/10_copy-device-udev-rules.conf
[Unit]
RequiresMountsFor=/system

[Service]
ExecStartPre=-/usr/lib/lxc-android/copy-udev-rules
